# deploy/supervisor_sample.ini
[program:api]
user=root
command=gunicorn -c deploy/prd_beedo.py beedo:app
autostart=true
autorestart=true
directory=/var/www/beedo
redirect_stderr=true
